# stage3/03-raspberry-tweaks/files/rc.local
#!/bin/sh -e

if grep -q 'Revision.*900093' /proc/cpuinfo; then
	ifup usb0
else
	ifup br-usb
	if grep -q dwc2 /boot/config.txt; then
		sed -i '/dwc2/d' /boot/config.txt
		reboot
	fi
fi
